# roles/common/templates/glances.service.j2
# Auto-generated by Ansible - manual changes will be lost on the next deploy.
[Unit]
Description=Glances web server (host stats for Home Assistant)
After=network.target

[Service]
# become: true (playbook.yml) defaults to root, so pipx installed this under
# /root/.local/bin - not a floating PATH lookup, since systemd services don't
# inherit an interactive shell's PATH.
ExecStart=/root/.local/bin/glances -w --port 61208
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
