# roles/minecraft/templates/mc-stats.service.j2
# Auto-generated by Ansible - manual changes will be lost on the next deploy.
[Unit]
Description=Minecraft stats JSON server (for Home Assistant)
After=network.target

[Service]
ExecStart=/usr/bin/python3 {{ mc_data_dir }}/mc_stats_server.py \
  --env-file /etc/minecraft/ha-stats.env \
  --data-dir {{ mc_data_dir }}/data \
  --port {{ mc_stats_port }}
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
